# dbus/bootkitd.service
[Unit]
Description=bookit daemon

[Service]
Type=dbus
BusName=org.opensuse.bootkit
ExecStart=/usr/sbin/bootkitd
